# .config/autostart/kmonad.desktop
[Desktop Entry]
    Name=kmonad
    Exec=kmonad /home/lam/.config/kmonad/config.kbd
    Type=Application
    X-GNOME-Autostart-enabled=true
